# config/tomcat.conf
[program:tomcat]
command=/opt/tomcat/bin/supervisor-wrapper.sh
